# packages/nvidia-k8s-device-plugin/nvidia-mps-control-daemon.service
[Unit]
Description=NVIDIA MPS Control Daemon
Before=nvidia-k8s-device-plugin.service

[Service]
ExecStart=/usr/bin/true
Restart=on-failure
RestartSec=2
RemainAfterExit=true

[Install]
RequiredBy=nvidia-k8s-device-plugin.service
