# packages/release/activate-preconfigured.service
[Unit]
Description=Activate preconfigured.target
DefaultDependencies=no
After=fipscheck.target
Requires=fipscheck.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl set-default preconfigured.target
ExecStart=/usr/bin/systemctl start preconfigured.target --no-block
RemainAfterExit=true
StandardOutput=tty
StandardError=inherit

[Install]
WantedBy=fipscheck.target
