# packages/release/activate-multi-user.service
[Unit]
Description=Activate multi-user.target
After=configured.target reboot-if-required.service
Requires=configured.target reboot-if-required.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl set-default multi-user.target
ExecStart=/usr/bin/systemctl start multi-user.target --no-block
RemainAfterExit=true
StandardError=journal+console

[Install]
WantedBy=configured.target
