# packages/release/activate-configured.service
[Unit]
Description=Activate configured.target
After=preconfigured.target
Requires=preconfigured.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl set-default configured.target
ExecStart=/usr/bin/systemctl start configured.target --no-block
RemainAfterExit=true
StandardError=journal+console

[Install]
WantedBy=preconfigured.target
