# packages/release/drivers.target
[Unit]
Description=Driver units
AllowIsolate=yes
After=basic.target
Before=preconfigured.target settings-applier.service
Requires=basic.target

[Install]
RequiredBy=preconfigured.target multi-user.target
