# packages/release/prepare-swap.service
[Unit]
Description=Prepare Swap
RefuseManualStart=true
RefuseManualStop=true
DefaultDependencies=no
Wants=modprobe@zram.service
After=modprobe@zram.service dev-zram0.device

[Service]
Type=oneshot
ExecStart=/usr/sbin/zramctl /dev/zram0 --size 1G
ExecStart=/usr/sbin/mkswap --uuid clear /dev/zram0
ExecStart=/usr/sbin/swapon --discard --priority 1024 /dev/zram0
RemainAfterExit=true
StandardError=journal+console
SyslogIdentifier=prepare-swap

[Install]
WantedBy=sysinit.target
