# packages/release/var-lib-bottlerocket.mount
[Unit]
Description=Private Directory (/var/lib/bottlerocket)
DefaultDependencies=no
Conflicts=umount.target
RequiresMountsFor=/var /.bottlerocket
Before=local-fs.target umount.target

[Mount]
What=/.bottlerocket
Where=/var/lib/bottlerocket
Type=none
Options=rbind,rshared

[Install]
WantedBy=preconfigured.target
