# packages/soci-snapshotter/soci-snapshotter.service
[Unit]
Description=soci snapshotter containerd plugin
Documentation=https://github.com/awslabs/soci-snapshotter
After=configured.target
Wants=configured.target
Before=containerd.service

[Service]
Type=notify
EnvironmentFile=/etc/network/proxy.env
ExecStart=/usr/bin/soci-snapshotter-grpc --address fd:// --config /etc/soci-snapshotter/config.toml --root /var/lib/soci-snapshotter
Restart=always
RestartSec=5
SyslogIdentifier=soci-snapshotter
